# packaging/thirdflare-one.service
[Unit]
Description=ThirdFlare One daemon
Documentation=https://github.com/oldrepublicwizard/thirdflare-one
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/usr/lib/thirdflare
EnvironmentFile=-/etc/default/thirdflare
Environment=THIRDFLARE_WEBUI=0
ExecStart=/usr/bin/env node /usr/lib/thirdflare/server.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
